# backend/dwd/crontab
 3-58/5 * * * * sh -c 'cd /usr/src/app && make update api OUT=/data/' > /proc/1/fd/1 2>/proc/1/fd/2
 @reboot sh -c 'cd /usr/src/app && make update api OUT=/data/' > /proc/1/fd/1 2>/proc/1/fd/2

