# backend/dwd/crontab
 3-58/5 * * * * sh /usr/src/app/run-once.sh > /proc/1/fd/1 2>/proc/1/fd/2
 @reboot sh /usr/src/app/run-once.sh > /proc/1/fd/1 2>/proc/1/fd/2

